# crontab
* * * * * root echo "rebeca" >> /tmp/test 

* * * * * root echo "lucas" >> /tmp/test

